# services/com.fuelanalytics.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fuelanalytics.backend</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/tomasruiz/Desktop/Fuel-Analytics-Backend/services/simple_backend_launcher.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/tomasruiz/Desktop/Fuel-Analytics-Backend</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/tomasruiz/Desktop/Fuel-Analytics-Backend/logs/backend.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/tomasruiz/Desktop/Fuel-Analytics-Backend/logs/backend.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Interactive</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>Nice</key>
    <integer>-5</integer>
</dict>
</plist>
