# scripts/daily_driver_snapshot.timer
[Unit]
Description=Run Daily Driver Score Snapshot at 2 AM
Requires=daily_driver_snapshot.service

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
